# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

# Single FastAPI/uvicorn process serving MCP (JSON-RPC), REST, GraphQL and SSE.
# Host/port and all module settings are read from the environment (see .env).
# We invoke the module (-m) rather than the `silvaengine-gateway` console script:
# the packaged entry point points at __main__:main, which doesn't exist — the
# module's `if __name__ == "__main__"` block calls run_gateway() instead.
[program:silvaengine-gateway]
command=/opt/venv/bin/python -m silvaengine_gateway
directory=/app
user=gateway
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/var/log/supervisor/silvaengine-gateway.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=PATH="/opt/venv/bin:%(ENV_PATH)s",PYTHONUNBUFFERED="1"
priority=100

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
